# roles/backups/templates/backup.service.j2
[Unit]
Description=Run Ansible backup playbook

[Service]
Type=oneshot
ExecStart=/usr/bin/ansible-playbook /etc/ansible/playbooks/backup.yml -e env=prod
WorkingDirectory=/etc/ansible