# roles/backups/templates/backup.service.j2
[Unit]
Description=Database and file backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh